# distros/C/CPAN-Testers-API/etc/runit/api/run
#!/bin/bash
. /opt/local/perlbrew/etc/bashrc
eval $( perl -Mlocal::lib )

export MOJO_HOME=$(pwd)

exec cpantesters-api prefork -l http://*:8000 2>&1
