# distros/C/CPAN-Testers-API/etc/runit/api/run
#!/bin/bash
. /opt/local/perlbrew/etc/bashrc
eval $( perl -Mlocal::lib )

export MOJO_HOME=$(pwd)
export MOJO_LOG_SHORT=1 # Remove timestamps from Mojo, added by log svc
export MOJO_INACTIVITY_TIMEOUT=300
export MOJO_HEARTBEAT_INTERVAL=30
export MOJO_HEARTBEAT_TIMEOUT=300

exec cpantesters-api prefork -I $MOJO_HEARTBEAT_INTERVAL -H $MOJO_HEARTBEAT_TIMEOUT -l http://*:8000 2>&1
