# distros/C/CPAN-Testers-API/etc/runit/api/run
#!/bin/bash
. /opt/local/perlbrew/etc/bashrc
eval $( perl -Mlocal::lib )

export MOJO_HOME=$(pwd)
export MOJO_LOG_SHORT=1 # Remove timestamps from Mojo, added by log svc

exec cpantesters-api prefork -l http://*:8000 2>&1
