# distros/C/CPAN-Testers-API/etc/runit/broker/run
#!/bin/bash
. /opt/local/perlbrew/etc/bashrc
eval $( perl -Mlocal::lib )
exec mercury broker -i 14400 -l http://*:9000 2>&1
