# distros/D/Daemon-Shutdown/examples/systemd/daemon-shutdown.service
[Unit]
Description = Daemon which shuts down box after certain triggers

[Service]
ExecStart = /usr/bin/sdd

[Install]
WantedBy = multi-user.target

