# distros/A/App-Padadoy/shares/app_padadoy/upstart.conf
description "foobare"

start on filesystem
stop on shutdown

respawn # restart, if ended unexpectedly

script
    chdir /home/foobar # base
    su -s /bin/sh -c 'padadoy start' foobar
end script
