# distros/C/CPAN-Testers-API/etc/runit/legacy-metabase/run
#!/bin/bash
. /opt/local/perlbrew/etc/bashrc
eval $( perl -Mlocal::lib )

export MOJO_HOME=$(pwd)

exec cpantesters-legacy-metabase prefork -l http://*:8250 2>&1
