# distros/C/CPAN-Testers-API/etc/runit/legacy-metabase/run
#!/bin/bash
. /opt/local/perlbrew/etc/bashrc
eval $( perl -Mlocal::lib )

export MOJO_HOME=$(pwd)

exec cpantesters-legacy-metabase prefork -l http://*:8250 -G 90 -H 90 -i 90  2>&1
