# distros/C/Config-Model-Systemd/t/model_tests.d/systemd-user-examples/basic-socket/gmail-imap-tunnel.socket
[Unit]
Description=Socket for Gmail IMAP tunnel

[Install]
WantedBy=sockets.target

[Socket]
ListenStream=9995
Accept=1
