# distros/m/meon-Web/etc/cron.d/meon-web
# crontab for meon::Web

40 3 * * * meonweb meon-web-expire-members
