# distros/L/Lemonldap-NG-Portal/scripts/lemonldap-ng-portal.timer
[Unit]
Description=Clean Lemonldap::NG sessions DB every 1h

[Timer]
OnCalendar=*-*-* *:07:07
Persistent=true

[Install]
WantedBy=timers.target
