# distros/C/Config-Model-Systemd/wr_root/cme-function/test.service
## This file was written by cme command.
## You can run 'cme edit systemd-service-file' to modify this file.
## You may also modify the content of this file with your favorite editor.


[Unit]
Description = test single unit
