# distros/N/Narada/skel/config/crontab
*/1  * * * *    exec chpst -L var/.lock.service runsvdir service/ &>/dev/null &
