# distros/W/WebService-Whistle-Pet-Tracker-API/systemd/perl-WebService-Whistle-Pet-Tracker-API-mqtt.service.d/override.conf
[Service]
Environment=MQTT_HOST=mqtt
Environment=WHISTLE_EMAIL=
Environment=WHISTLE_PASSWORD=
