# distros/C/Config-Model-Systemd/t/model_tests.d/systemd-user-examples/duplicity/duplicity.timer
[Timer]
OnCalendar = 13:00

[Unit]
Description = Run duplicity
