# distros/N/Narada/skel/service/log/run
#!/bin/sh
exec socklog unix ../../var/log.sock 2>&1
