# distros/B/Bot-BasicBot-Pluggable-Module-Crontab/t/data/01test.cron
# a comment

* * * * * * #dev Minute Check!

