# distros/F/Fwctl/debian/postinst
#! /bin/sh
set -e

if [ "$1" = "configure" ]
then
  update-rc.d fwctl defaults >/dev/null
fi

