# distros/L/Lemonldap-NG-Handler/scripts/lemonldap-ng-handler.timer
[Unit]
Description=Purge Lemonldap::NG handler cache every 30 minutes

[Timer]
OnCalendar=*-*-* *:17,47:00
Persistent=true

[Install]
WantedBy=timers.target
